# crontab/delete-log-files
#DMZ & MID trong /app/logs
* * * * * root find #LOG_DIR# -type f -mtime #MAX_LOGS_DAY# -delete

